# config/99-dvd-ripper.rules
# Udev rule for automatic DVD ripping
# Triggers systemd service when DVD media is inserted
# Uses systemd-run to escape udev's cgroup/environment restrictions
#
# ID_CDROM_MEDIA=1 ensures we only trigger on insert, not eject

KERNEL=="sr[0-9]*", SUBSYSTEM=="block", ENV{ID_CDROM}=="1", ACTION=="change", ENV{ID_CDROM_MEDIA}=="1", RUN+="/usr/bin/systemd-run --no-block /usr/local/bin/dvd-ripper.sh /dev/%k"
