# config/dvd-ripper@.service
[Unit]
Description=DVD Ripper for %I (Legacy Mode)
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/local/bin/dvd-ripper.sh /dev/%i
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dvd-ripper

# User isolation - run as dvd-rip user (needs DVD device access)
User=dvd-rip
Group=dvd-ripper
SupplementaryGroups=cdrom

# libdvdcss cache directory (user has no home dir)
Environment=DVDCSS_CACHE=/var/cache/dvdcss
Environment=HOME=/var/cache/dvdcss

# Runtime directory for lock files (created automatically)
# Mode 0775 allows all dvd-ripper group members to write lock files
RuntimeDirectory=dvd-ripper
RuntimeDirectoryMode=0775
RuntimeDirectoryPreserve=yes

# Resource limits (prevent runaway processes)
CPUQuota=80%
MemoryMax=2G
TasksMax=32
IOWeight=50
Nice=10

# Timeout - give it plenty of time (4 hours for long DVDs)
TimeoutStartSec=0
TimeoutStopSec=30

# Don't restart on failure
Restart=no

# Security hardening
# Note: Using 'true' instead of 'strict' to avoid mount namespace issues
# Note: DVD access requires some capabilities for eject/ioctl operations
# Note: PrivateTmp disabled - staging dir must be shared between pipeline stages
ProtectSystem=true
ProtectHome=true
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
